# backup.timer
[Unit]
Description=Timer for backup

[Timer]
OnBootSec=1min
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target